# scripts/com.qkangber.signal-log.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  每日訊號快照排程（macOS launchd 版）。範本，路徑寫死本機專案位置——換機器要改 Program/WorkingDirectory。
  安裝：cp scripts/com.qkangber.signal-log.plist ~/Library/LaunchAgents/
        launchctl bootstrap gui/$UID ~/Library/LaunchAgents/com.qkangber.signal-log.plist
  手動觸發一次：launchctl kickstart -p gui/$UID/com.qkangber.signal-log
  移除：launchctl bootout gui/$UID/com.qkangber.signal-log
  睡眠中錯過的時段，launchd 會在喚醒後補跑一次（快照本身也會把資料補進正確的資料日檔案）。
-->
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.qkangber.signal-log</string>

  <key>Program</key>
  <string>/Users/kc/qkangber-toolbox/scripts/snapshot-task.sh</string>

  <key>WorkingDirectory</key>
  <string>/Users/kc/qkangber-toolbox</string>

  <!-- 平日 14:30 收盤後跑；16:00 再補一次，網路剛醒/腳本三輪重試仍失敗時還有第二次機會 -->
  <key>StartCalendarInterval</key>
  <array>
    <dict><key>Weekday</key><integer>1</integer><key>Hour</key><integer>14</integer><key>Minute</key><integer>30</integer></dict>
    <dict><key>Weekday</key><integer>2</integer><key>Hour</key><integer>14</integer><key>Minute</key><integer>30</integer></dict>
    <dict><key>Weekday</key><integer>3</integer><key>Hour</key><integer>14</integer><key>Minute</key><integer>30</integer></dict>
    <dict><key>Weekday</key><integer>4</integer><key>Hour</key><integer>14</integer><key>Minute</key><integer>30</integer></dict>
    <dict><key>Weekday</key><integer>5</integer><key>Hour</key><integer>14</integer><key>Minute</key><integer>30</integer></dict>
    <dict><key>Weekday</key><integer>1</integer><key>Hour</key><integer>16</integer><key>Minute</key><integer>0</integer></dict>
    <dict><key>Weekday</key><integer>2</integer><key>Hour</key><integer>16</integer><key>Minute</key><integer>0</integer></dict>
    <dict><key>Weekday</key><integer>3</integer><key>Hour</key><integer>16</integer><key>Minute</key><integer>0</integer></dict>
    <dict><key>Weekday</key><integer>4</integer><key>Hour</key><integer>16</integer><key>Minute</key><integer>0</integer></dict>
    <dict><key>Weekday</key><integer>5</integer><key>Hour</key><integer>16</integer><key>Minute</key><integer>0</integer></dict>
  </array>

  <!-- 腳本自己會寫 run.log，這裡只留 launchd 層級（連腳本都沒起來）的錯誤 -->
  <key>StandardErrorPath</key>
  <string>/Users/kc/qkangber-toolbox/data/signal-log/launchd.err</string>
</dict>
</plist>
